# systemd/system/certbot-renew.service.d/local.conf
[Service]
ExecStart=/usr/bin/certbot -q renew --no-random-sleep-on-renew
