# systemd/system/fstrim.service.d/local.conf
[Service]
CPUSchedulingPolicy=idle
IOSchedulingClass=idle
