# systemd/system/fstrim.service.d/override.conf
[Unit]
Wants=xfs_fsr.service
After=xfs_fsr.service

[Service]
CPUSchedulingPolicy=idle
IOSchedulingClass=idle
