# etc/systemd/system/mdmonitor.service.d/override.conf
[Service]
ExecStart=
ExecStart=/sbin/mdadm --monitor --scan --syslog
