# systemd/system/systemd-journald-vacuum.timer
[Unit]
Description=Run systemd-journald-vacuum hourly

[Timer]
OnCalendar=hourly

[Install]
WantedBy=timers.target
