# etc/systemd/system/sshd.service.d/override.conf
[Service]
LimitNOFILE=8192
ManagedOOMPreference=avoid
RestartMaxDelaySec=10s
RestartSec=100ms
RestartSteps=5
