# etc/systemd/system/sshd.service.d/override.conf
[Unit]
After=network-online.target

[Service]
LimitNOFILE=8192
ManagedOOMPreference=avoid
