# etc/systemd/system/fstrim.timer.d/override.conf
[Unit]
Description=Discard unused filesystem blocks once a day

[Timer]
OnCalendar=
OnCalendar=daily
