# etc/systemd/system/systemd-boot-update.service.d/local.conf
[Service]
Environment=SYSTEMD_RELAX_ESP_CHECKS=1
