# etc/systemd/system/systemd-cryptsetup@.service.d/override.conf
[Service]
ExecStartPost=/usr/bin/udevadm trigger /dev/mapper/%i
