# systemd/system/sysstat-collect.timer.d/local.conf
[Unit]
Description=Run system activity accounting tool every minute

[Timer]
AccuracySec=1us
OnCalendar=
OnCalendar=minutely
