# systemd/system/sysstat-collect.timer.d/minutely.conf
[Unit]
Description=Run system activity accounting tool every minute

[Timer]
OnCalendar=
OnCalendar=minutely
