# etc/systemd/system/syslog-ng@default.service.d/override.conf
[Service]
Restart=always
RestartMaxDelaySec=10s
RestartSec=100ms
RestartSteps=5
