# systemd/system/attestation.service.d/local.conf
[Service]
IPAddressAllow={{ipv4_address}}
IPAddressAllow={{ipv6_address}}
