# etc/systemd/system/rotate-session-ticket-keys.timer
[Unit]
Description=Run rotate-session-ticket-keys every 8 hours
After=create-session-ticket-keys.service
Requires=create-session-ticket-keys.service

[Timer]
OnActiveSec=8h
OnUnitActiveSec=8h

[Install]
WantedBy=timers.target
