# etc/systemd/system/rotate-session-ticket-keys.timer
[Unit]
Description=Run rotate-session-ticket-keys three times daily

[Timer]
OnCalendar=0/8:00:00

[Install]
WantedBy=timers.target
