# etc/systemd/system/tcp-fastopen-rotate-keys.service
[Unit]
Description=Rotate TCP Fast Open keys

[Service]
ExecStart=/usr/local/bin/tcp-fastopen-rotate-keys
Type=oneshot
UMask=0077
